# monitoring-test.service
[Unit]
Description=Test Process Monitoring Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring_test.sh
User=root
